# moderator/systemd/schema-moderator.service
[Unit]
Description=Renaissance Schema Moderator Service (BUILD-SPEC-v2) — DDL review authority + approval ledger
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=/opt/moderator/bin
EnvironmentFile=/opt/moderator/moderator.env
ExecStart=/opt/duckdb/venv/bin/python /opt/moderator/bin/moderator_server.py
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
